# etc/rc.local
/usr/local/sbin/dhcp6c EXTIF
/usr/sbin/rtadvd INTIF
